# docAndTips/dockersample/10.85.47.54/z/conf/snap8030.conf
[supervisord]
nodaemon=true


[program:xvfb]
command=bash -c "/usr/bin/Xvfb :0 -screen 0 1024x768x24 &"

[program:snap]
command=java -cp classes:/usr/local/joblib/* -Xms256m -Xmx1024m -Xmn128m com.sina.z.job.ZSnapServer 0.0.0.0 8030 8031 /usr/local/tomcat/snap.log &

